# src/crontab
SHELL=/bin/bash
BASH_ENV=/container.env
*/10 * * * * /scripts/exampleScript.sh >> /var/log/cron.log 2>&1
# dont remove this comment or this will no longer be a valid cron file.